# roles/common/files/ansible-pull/ansible-pull.timer
[Unit]
Description=Update config with ansible daily.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
